# rootfs/etc/services.d/nginx/run
#!/bin/sh
exec nginx -g "pid /run/nginx.pid; daemon off;";